# RaspberryPi/rc.local
# /etc/rc.local file for homepass on Raspbian
# HeadURL: https://github.com/danielhoherd/homepass/blob/master/RaspberryPi/rc.local

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi
service hostapd stop
screen -S HomePassScript bash -c "/root/homepass.sh | tee -a /root/homepass.log"
exit 0
